# infrastructure/jupiter/ansible/files/opencloud-backup.service
[Unit]
Description=OpenCloud Backup to Scaleway Object Storage
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/opencloud-backup.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencloud-backup

[Install]
WantedBy=multi-user.target
